# backend/avatar-bot.service
[Unit]
Description=AVATAR Matrix Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/avatar-new/backend
ExecStart=/root/avatar-new/backend/venv/bin/python bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=avatar-bot

[Install]
WantedBy=multi-user.target
