# backend/avatar-backend.service
[Unit]
Description=AVATAR Matrix Backend (FastAPI + uvicorn)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/avatar-new/backend
ExecStart=/root/avatar-new/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=avatar-new-backend

[Install]
WantedBy=multi-user.target
